# sigsync_recv.ini
[program:sigsync-recv]
command=nohup /opt/sigsync/sigsync daemon
directory=/opt/sigsync
user=root
stdout_logfile=/tmp/stdout_logfile.log
stdout_logfile_backups=10
stdout_events_enabled=false
stderr_logfile=/tmp/stderr_logfile.log
stderr_logfile_maxbytes=1MB
stderr_logfile_backups=10
stderr_events_enabled=false
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true
